# Spielehub-V1/Schiffe-Versenken-main/battleship-deluxe.service
[Unit]
Description=Battleship Deluxe (Node WebSocket Backend)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/battleship-deluxe
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8080
EnvironmentFile=-/opt/battleship-deluxe/.env
ExecStart=/usr/bin/node /opt/battleship-deluxe/server.js
Restart=always
RestartSec=3
User=www-data
Group=www-data

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
